# files/systemd/i3-session.target
[Unit]
Description=i3 session running
BindsTo=graphical-session.target
After=i3-session-polybar.target
Before=i3-session-tray.target
