# files/systemd/i3-session-tray.target
[Unit]
Description=Launch tray apps
BindsTo=graphical-session.target
After=i3-session.target
