# files/systemd/i3-session-pre.target
[Unit]
Description=i3 session setup
BindsTo=graphical-session-pre.target
Before=i3-session-polybar.target
